# Actividades/Actividad 3/saludo.service
[Unit]
Description=Servicio saludo

[Service]
ExecStart=/home/usuario/script.sh

[Install]
WantedBy=multi-user.target